# conf/roles/slurm/templates/slurmrestd.service.j2
[Service]
# Unset vendor unit ExecStart to avoid cumulative definition
ExecStart=
{% if slurm_with_jwt %}
Environment="SLURM_JWT=daemon"
ExecStart=/usr/sbin/slurmrestd $SLURMRESTD_OPTIONS [::]:6820
{% else %}
Environment=
Environment=SLURMRESTD_SECURITY=disable_user_check
ExecStart=/usr/sbin/slurmrestd $SLURMRESTD_OPTIONS unix:{{ slurm_restd_socket }}
RuntimeDirectory=slurmrestd
RuntimeDirectoryMode=0755
{% endif %}
User=slurm
Group=slurm
